# conf/handset-daemon.service
[Unit]
Description=HP Internet Handset user-space driver daemon
Documentation=https://github.com/arkygeek/voipHandset
After=graphical-session.target pipewire.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/handset daemon
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=graphical-session.target
